# arkcase.ini
[program:arkcase]
command=/app/tomcat/bin/catalina.sh run
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
user=core
autorestart=false
directory=/app/home
